# _cgroup2_find_path
#!/sbin/openrc-run
. /lib/rc/sh/rc-cgroup.sh
cgroup2_find_path
