# _cgroup2_find_path
#!/sbin/openrc-run
. /usr/libexec/rc/sh/rc-cgroup.sh
cgroup2_find_path
